# scripts/systemd/machtblick-bundestag-auto-refresh.service.in
[Unit]
Description=Machtblick Bundestag auto refresh

[Service]
Type=oneshot
WorkingDirectory=__ROOT__
Environment=MACHTBLICK_ROOT=__ROOT__
ExecStart=__ROOT__/scripts/scheduled-bundestag-auto-refresh
